# infra/launchd/com.openclaw.clawndom.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.clawndom</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>__INSTALL_PATH__/dist/server.js</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>OPENCLAW_TOKEN</key>
    <string>__OPENCLAW_TOKEN__</string>
    <key>PROVIDERS_CONFIG</key>
    <string>__PROVIDERS_CONFIG__</string>
    <key>OPENCLAW_GATEWAY_WS_URL</key>
    <string>ws://127.0.0.1:18789</string>
    <key>REDIS_URL</key>
    <string>__REDIS_URL__</string>
    <key>PORT</key>
    <string>__PORT__</string>
    <key>MAX_CONCURRENT_RUNS</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/usr/local/var/log/clawndom.log</string>

  <key>StandardErrorPath</key>
  <string>/usr/local/var/log/clawndom.log</string>
</dict>
</plist>
